# systemd/user/gpg-agent.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gpg-agent --daemon --options /home/${USER}/.gnupg/gpg-agent.conf
